# internal/service/testdata/systemd_golden_noenv.service
[Unit]
Description=AgentPaaS Control Daemon

[Service]
ExecStart=/opt/agentpaas/bin/agentpaasd --home /home/testuser/.agentpaas
Restart=on-failure

[Install]
WantedBy=default.target
